# config/crontab
# AviationWX Cron Jobs
# Runs inside Docker container

# Set PATH for cron jobs (cron has minimal PATH by default)
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Disable mail output (send to /dev/null) - we want output in Docker logs instead
MAILTO=""

# Webcam Image Fetcher - Runs every minute
# Fetches webcam images from configured sources (MJPEG streams, RTSP, static images)
# and caches them for serving via the webcam API. Skips cameras with fresh cache.
# Output and errors are logged to Docker logs (stdout/stderr).
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/fetch-webcam.php 2>&1

# Weather Data Fetcher - Runs every minute
# Refreshes weather cache for all configured airports by calling the weather API endpoint.
# Ensures weather data stays fresh for display on airport pages.
# Output and errors are logged to Docker logs (stdout/stderr).
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/fetch-weather.php 2>&1

# Push Webcam Processor - Runs every minute
# Processes webcam images uploaded via SFTP/FTP/FTPS by push cameras.
# Moves uploaded images from uploads directory to cache directory for serving.
# Output and errors are logged to Docker logs (stdout/stderr).
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/process-push-webcams.php 2>&1

# Let's Encrypt Certificate Renewal - Runs weekly (Sunday at 3 AM)
# Renews SSL certificates for upload.aviationwx.org if they're expiring soon.
# After renewal, enables SSL/TLS in vsftpd if certificates are available.
# Errors are logged but don't fail the container (|| true).
0 3 * * 0 root certbot renew --quiet --no-self-upgrade && /usr/local/bin/enable-vsftpd-ssl.sh || true
