# config/crontab
# AviationWX Cron Jobs
# Runs inside Docker container
# Format: minute hour day month weekday user command

# Set PATH for cron jobs (cron has minimal PATH by default)
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Webcam refresh - every minute
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/fetch-webcam.php > /dev/null 2>&1

# Weather refresh - every minute
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/fetch-weather.php > /dev/null 2>&1

# Push webcam file processor - every minute
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/process-push-webcams.php > /dev/null 2>&1

# Push webcam config sync - every minute
# Must run as root to write to /etc/vsftpd/ and create system users
# Log errors to file for debugging (stdout/stderr go to Docker logs via logger)
# PATH is set above to ensure db_load, db_verify, etc. are found
* * * * * root cd /var/www/html && /usr/local/bin/php scripts/sync-push-config.php >> /var/log/cron-sync-push-config.log 2>&1

# Let's Encrypt certificate renewal - weekly (Sunday at 3 AM)
# After renewal, enable SSL in vsftpd if not already enabled
0 3 * * 0 root certbot renew --quiet --no-self-upgrade && /usr/local/bin/enable-vsftpd-ssl.sh > /dev/null 2>&1 || true

