# config/crontab
# AviationWX Cron Jobs
# Runs inside Docker container

# Set PATH for cron jobs (cron has minimal PATH by default)
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Disable mail output (send to /dev/null) - we use file-based logging
MAILTO=""
# Set environment variables for cron jobs (cron doesn't inherit container env by default)
# CONFIG_PATH is critical for scripts to find airports.json
CONFIG_PATH=/var/www/html/secrets/airports.json

# Cron Heartbeat - Runs every minute
# Logs a heartbeat message to confirm cron daemon is running and executing jobs
# This helps monitor cron health and detect when cron stops working
# Logs to /var/log/aviationwx/ (logrotate handles rotation: 7 days, 100MB max)
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/cron-heartbeat.php >> /var/log/aviationwx/cron-heartbeat.log 2>&1

# Scheduler Health Check - Runs every minute
# Validates scheduler daemon is running and healthy, restarts if needed.
# The scheduler daemon handles all weather and webcam updates (replaces cron-based fetching).
# Logs to /var/log/aviationwx/ via aviationwx_log() function.
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/scheduler-health-check.php >> /var/log/aviationwx/scheduler-health-check.log 2>&1

# Memory Sampler - Runs every minute (samples 12x at 5s intervals internally)
# Collects memory usage samples for CPU-load-style rolling averages (1m/5m/15m).
# Samples are stored in hourly JSON buckets following metrics system conventions.
# Logs to /var/log/aviationwx/ via aviationwx_log() function.
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/sample-memory.php >> /var/log/aviationwx/memory-sampler.log 2>&1

# NOTE: Push webcam processing is now handled by the unified scheduler daemon.
# The scheduler uses WebcamScheduleQueue for efficient scheduling of ALL webcams
# (both push and pull) with per-camera refresh_seconds configuration.
# This cron job has been removed as part of the unified webcam worker refactor.

# Cache Cleanup - Runs daily at 4 AM
# Two-layer cleanup strategy (belt and suspenders):
# - Layer 1: Primary cleanup for files without automatic cleanup (rate limits, error files)
# - Layer 2: Backup cleanup as safety net for files that should have auto-cleanup
# - Layer 3: Orphan cleanup for removed airports
# Logs to /var/log/aviationwx/ via aviationwx_log() function.
0 4 * * * www-data cd /var/www/html && /usr/local/bin/php scripts/cleanup-cache.php >> /var/log/aviationwx/cleanup-cache.log 2>&1

# NOTE: Certificate renewal is handled by the HOST's systemd certbot.timer (runs twice daily)
# The deploy hook at /etc/letsencrypt/renewal-hooks/deploy/aviationwx-deploy-hook.sh
# automatically copies renewed certificates and reloads services.
# See scripts/certbot-deploy-hook.sh for the hook implementation.
