# config/crontab
# AviationWX Cron Jobs
# Runs inside Docker container
# Format: minute hour day month weekday user command

# Webcam refresh - every minute
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/fetch-webcam.php > /dev/null 2>&1

# Weather refresh - every minute
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/fetch-weather.php > /dev/null 2>&1

# Push webcam file processor - every minute
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/process-push-webcams.php > /dev/null 2>&1

# Push webcam config sync - every minute
* * * * * www-data cd /var/www/html && /usr/local/bin/php scripts/sync-push-config.php > /dev/null 2>&1

# Let's Encrypt certificate renewal - weekly (Sunday at 3 AM)
0 3 * * 0 root certbot renew --quiet --no-self-upgrade > /dev/null 2>&1

